# infrastructure-nifi/playbooks/registry/nifi-registry.timer
[Unit]
Description=Starts NiFi after boot.

[Timer]
OnBootSec=1min
Unit=nifi-registry.service

[Install]
WantedBy=multi-user.target
